# scripts/observability/templates/systemd/phpfpm_exporter.service
[Unit]
Description=Prometheus PHP-FPM Exporter
Documentation=https://github.com/hipages/php-fpm_exporter
After=network-online.target php-fpm.service
Wants=network-online.target

[Service]
Type=simple
User=prometheus
Group=prometheus
ExecStart=/usr/local/bin/phpfpm_exporter \
    --web.listen-address=0.0.0.0:9253 \
    --phpfpm.scrape-uri=tcp://localhost:9000/status

# Alternative for Unix socket:
# ExecStart=/usr/local/bin/phpfpm_exporter \
#     --web.listen-address=0.0.0.0:9253 \
#     --phpfpm.scrape-uri=unix:///var/run/php-fpm/www.sock;status

SyslogIdentifier=phpfpm_exporter
Restart=always
RestartSec=10s

# Security hardening
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true

# Resource limits
LimitNOFILE=8192
LimitNPROC=512

[Install]
WantedBy=multi-user.target
