# docker/web/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info

# ============================================================================
# Core Services
# ============================================================================

[program:mysql]
command=/usr/sbin/mysqld --defaults-file=/etc/mysql/my.cnf
user=mysql
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/mysql.log
stderr_logfile=/var/log/supervisor/mysql_err.log
priority=10

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
user=redis
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis_err.log
priority=20

[program:php-fpm]
command=/usr/sbin/php-fpm8.4 --nodaemonize --fpm-config /etc/php/8.4/fpm/php-fpm.conf
user=root
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/php-fpm.log
stderr_logfile=/var/log/supervisor/php-fpm_err.log
priority=30

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
user=root
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log
priority=40

# ============================================================================
# Laravel Application Services
# ============================================================================

[program:laravel-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/chom/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --timeout=300
user=www-data
numprocs=2
autostart=true
autorestart=true
startretries=10
stopwaitsecs=3600
stdout_logfile=/var/log/supervisor/laravel-queue.log
stderr_logfile=/var/log/supervisor/laravel-queue_err.log
priority=50

[program:laravel-scheduler]
command=/bin/bash -c "while true; do php /var/www/chom/artisan schedule:run --verbose --no-interaction; sleep 60; done"
user=www-data
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/laravel-scheduler.log
stderr_logfile=/var/log/supervisor/laravel-scheduler_err.log
priority=60

# ============================================================================
# Monitoring Exporters
# ============================================================================

[program:node-exporter]
command=/usr/local/bin/node_exporter \
    --collector.filesystem.mount-points-exclude='^/(sys|proc|dev|host|etc)($$|/)' \
    --collector.netclass.ignored-devices='^(veth.*|docker.*|br-.*)$$'
user=root
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/node-exporter.log
stderr_logfile=/var/log/supervisor/node-exporter_err.log
priority=70

[program:nginx-exporter]
command=/usr/local/bin/nginx-prometheus-exporter -nginx.scrape-uri=http://localhost:8080/stub_status
user=root
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/nginx-exporter.log
stderr_logfile=/var/log/supervisor/nginx-exporter_err.log
priority=71

[program:mysql-exporter]
command=/usr/local/bin/mysqld_exporter --config.my-cnf=/etc/mysql/exporter.cnf
user=root
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/mysql-exporter.log
stderr_logfile=/var/log/supervisor/mysql-exporter_err.log
priority=72
environment=DATA_SOURCE_NAME="exporter:exporter_password@(localhost:3306)/"

[program:php-fpm-exporter]
command=/usr/local/bin/php-fpm_exporter --phpfpm.scrape-uri=unix:///run/php/php8.4-fpm.sock;/fpm-status
user=root
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/php-fpm-exporter.log
stderr_logfile=/var/log/supervisor/php-fpm-exporter_err.log
priority=73

# ============================================================================
# Grafana Alloy (Metrics & Logs Shipper)
# ============================================================================

[program:alloy]
command=/usr/local/bin/alloy run /etc/alloy/config.alloy \
    --storage.path=/var/lib/alloy \
    --server.http.listen-addr=0.0.0.0:12345
user=root
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/alloy.log
stderr_logfile=/var/log/supervisor/alloy_err.log
priority=80

# ============================================================================
# Cron for Laravel Scheduler (alternative approach)
# ============================================================================

[program:cron]
command=/usr/sbin/cron -f
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/cron.log
stderr_logfile=/var/log/supervisor/cron_err.log
priority=90
